# root/etc/s6-overlay/s6-rc.d/crond/up
/etc/s6-overlay/scripts/setup_ipv6_route